# zomboid_supervisord.conf
[program:zomboid]
command=bash -c /home/jd/patch-zomboid.sh
autorestart=true
user=jd
